# deploy/local-cd/systemd/watchlist-deploy.timer
[Unit]
Description=Poll for CI-validated Watchlist movie releases

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=60s
Persistent=true
Unit=watchlist-deploy.service

[Install]
WantedBy=timers.target
