# deploy/local-cd/systemd/watchlist-deploy.service
[Unit]
Description=Deploy CI-validated Watchlist movie services
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=watchlist
Group=watchlist
UMask=0077
EnvironmentFile=/opt/watchlist-prod/config/deploy.env
Environment=HOME=/opt/watchlist-prod
Environment=DOCKER_CONFIG=/opt/watchlist-prod/.docker
ExecStart=/bin/bash /opt/watchlist-prod/deployer/deploy-movie-sync.sh
WorkingDirectory=/opt/watchlist-prod
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/opt/watchlist-prod
